# Compiti/Gestione_di_servizi_monitoraggio/mymon.service/mymon.service
[Unit]
Description=Servizio di monitoraggio
Requires=mylog.service

[Service]
Type=simple
ExecStart=mymon.service.sh
